# srcs/nginx/srcs/supervisor.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx
autostart=true
autorestart=true

[program:ssh]
command=php-fpm7
autostart=true
autorestart=true